# systemd/erpsync-gunicorn.service
# ERP Sync Gunicorn Service
# Copy to: /etc/systemd/system/erpsync-gunicorn.service

[Unit]
Description=ERP Sync Gunicorn WSGI Server
Documentation=https://docs.djangoproject.com/
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/erpsync

# Environment
Environment="PATH=/opt/erpsync/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/erpsync"
Environment="DJANGO_SETTINGS_MODULE=erpsync_project.settings"

# Gunicorn command
ExecStart=/opt/erpsync/venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --worker-class sync \
    --worker-connections 1000 \
    --max-requests 1000 \
    --max-requests-jitter 50 \
    --timeout 60 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --access-logfile /var/log/erpsync/gunicorn-access.log \
    --error-logfile /var/log/erpsync/gunicorn-error.log \
    --log-level info \
    --capture-output \
    --enable-stdio-inheritance \
    erpsync_project.wsgi:application

# Restart policy
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGQUIT
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
